# build/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
pidfile=/var/run/supervisord.pid
loglevel=info

[program:gcopy-backend]
command=/app/backend/start.sh
directory=/app/backend
autostart=true
autorestart=true
startretries=3
startsecs=5
# Output to stdout/stderr for Docker logs
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
# Supervisor automatically adds process name prefix to logs

[program:gcopy-frontend]
command=/bin/sh -c 'until curl -fs http://127.0.0.1:3376/api/v1/ping; do sleep 1; done; node server.js'
directory=/app/frontend
autostart=true
autorestart=true
startretries=3
startsecs=5
# Output to stdout/stderr for Docker logs
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
# Supervisor automatically adds process name prefix to logs
